# systemd/azerothcore-auth.service
# Generic acore-manager authserver service template.
# Copy this file to /etc/systemd/system/ and adjust User, Group, and paths if
# your local installation differs from the defaults shown here.
#
# Runtime configs are expected through the active release:
# /opt/acore-manager/current/etc/authserver.conf -> /opt/acore-manager/shared/configs/authserver.conf

[Unit]
Description=AzerothCore Auth Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=azerothcore
Group=azerothcore
WorkingDirectory=/opt/acore-manager/current/bin
ExecStart=/opt/acore-manager/current/bin/authserver
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
